# sensor8.service
[Unit]
Description=Sensor 8
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/viz-debug
Restart=on-failure
RestartSec=5s
ExecStart=cargo run server /dev/sensor8 9998

[Install]
WantedBy=multi-user.target
