# sensor5.service
[Unit]
Description=Sensor 5
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/viz-debug
Restart=on-failure
RestartSec=5s
ExecStart=cargo run server /dev/sensor5 9995

[Install]
WantedBy=multi-user.target
